# deploy/systemd/sentinel-analyze.service
[Unit]
Description=Sentinel analyze — run 7 detection passes + route findings
Documentation=https://github.com/chitinhq/sentinel
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/sentinel/env
Environment=SENTINEL_DIGEST_DIR=%h/.local/share/sentinel/digests
Environment=SENTINEL_CONFIG=%h/workspace/sentinel-heartbeat/sentinel.yaml
WorkingDirectory=%h/workspace/sentinel-heartbeat
ExecStartPre=/bin/mkdir -p %h/.local/share/sentinel/digests
ExecStart=%h/workspace/sentinel-heartbeat/bin/sentinel analyze
SuccessExitStatus=0
# Analyze can take longer; cap it so a stuck run doesn't block the next tick.
TimeoutStartSec=30min
Nice=10

[Install]
WantedBy=default.target
